# rootfs/etc/services.d/20-cloudplow/run
#!/usr/bin/with-contenv sh

echo run cloudplow

exec s6-setuidgid abc cloudplow --config /data/cloudplow/config.json --cachefile /data/cloudplow/cache.db --logfile /logs/cloudplow.log run

