# rootfs/etc/services.d/20-syslog-replication/run
#!/usr/bin/with-contenv sh

echo run syslog replication

exec s6-setuidgid abc tail -f /syslog >> /logs/unraid.log

