# rootfs/etc/services.d/10-cron/run
#!/command/with-contenv sh

rm -f /tmp/update_docker.lock
echo run crond

crond
exec s6-setuidgid abc /usr/sbin/crond -f
