# rootfs/etc/services.d/10-cron/run
#!/usr/bin/with-contenv sh

echo run crond

exec s6-setuidgid abc /usr/sbin/crond -f

